# services/akros-dns.service
# Copy (chmod 664) to /etc/systemd/system/akros-dns.service
#
# systemctl daemon-reload
# systemctl start akros-dns
# systemctl status akros-dns
# systemctl enable akros-dns
#
# podman pull coredns/coredns

[Unit]
Description=CoreDNS DNS Server
After=network.target syslog.target

[Service]
Type=simple
TimeoutStartSec=5m
ExecStartPre=-/usr/bin/podman rm "akros-dns"

ExecStart=/usr/bin/podman run --name akros-dns -p 53:53/udp -v /opt/coredns/Corefile:/Corefile:ro -v /opt/coredns/zones:/zones:ro coredns/coredns

ExecReload=-/usr/bin/podman stop "akros-dns"
ExecReload=-/usr/bin/podman rm "akros-dns"
ExecStop=-/usr/bin/podman stop "akros-dns"
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target

#### NOTES
# podman run --rm --name akros-dns -p 53:53 -v /opt/coredns/Corefile:/Corefile:ro coredns/coredns -conf /Corefile 
